# etc/system/usr/lib/systemd/system/ddate_update.service
[Unit]
Description=Daily ddate update for browser startpage

[Service]
Type=simple
ExecStart=/bin/sh /usr/local/bin/ddate_update
Restart=on-abnormal
