# etc/system/usr/lib/systemd/system/ddate_update.timer
[Unit]
Description=Daily ddate update for browser startpage

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
